# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:sshd]
command=/usr/sbin/sshd -D
